# blobstream-monitor/services/monitor_orchestrator.service
[Unit]
Description=Monitor Orchestrator Logs and Summarize in JSON
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor_orchestrator_json.sh
Restart=on-failure
RestartSec=5min

[Install]
WantedBy=multi-user.target
